# scripts/systemd/boxbot.service
[Unit]
Description=boxBot — ambient Claude agent
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=jhart
Group=jhart
WorkingDirectory=/home/jhart/software/boxBot

# Cap glibc malloc arenas. Default is 8*ncores (32 on the Pi 5), which
# fragments badly under concurrent mixed-size allocation across the
# voice/perception/agent subsystems — freed memory ends up stranded in
# per-arena free lists and never returns to the OS, building RSS across
# conversations until OOM. 2 arenas keeps fragmentation bounded; the
# in-process malloc_trim() (diagnostics.memory) returns the rest.
Environment=MALLOC_ARENA_MAX=2

ExecStart=/home/jhart/software/boxBot/.venv/bin/boxbot

# Restart on crash, but not when the process exits cleanly (deploy /
# explicit stop / RSS-guardrail shutdown all use SIGTERM → exit 0).
Restart=on-failure
RestartSec=5s

# Cgroup memory limits. MemoryHigh exerts pressure but doesn't kill;
# our in-process RSS guardrail (config.diagnostics.memory.rss_shutdown_mb)
# is expected to trip well before MemoryMax does. MemoryMax is the hard
# backstop so a runaway never takes down the whole box.
MemoryHigh=6G
MemoryMax=7G

# Surface stdout/stderr for crashes that happen before the Python logger
# is initialised (or that bypass it entirely — tracebacks on import,
# kernel OOM messages, etc.). Application logs continue to go to
# logs/boxbot.log via RotatingFileHandler in Python.
StandardOutput=append:/home/jhart/software/boxBot/logs/boxbot-systemd.log
StandardError=inherit

# Give the agent time to flush memory extraction, close DB connections,
# and stop sandbox runners on shutdown. The graceful path in main.py
# walks the subsystems in reverse order.
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
